# systemd/nuc16pro-scx-server.service
[Unit]
Description=sched_ext server scheduler - ASUS NUC 16 Pro (Panther Lake)
Documentation=https://github.com/sched-ext/scx
# Attach the scheduler before dockerd starts the container fleet. sched_ext runs
# ops.cgroup_init() once per existing cgroup at attach time; attaching after the
# ~20 docker-<app>.service units launch means ~175 cgroups are initialized in one
# batch under boot slab pressure, which transiently fails -ENOMEM (cgroup_init -12)
# and demotes the scheduler to the fallback. Attaching before docker.service inits
# only a handful of cgroups up front, then each container's cgroup is initialized
# incrementally as it starts. After=basic.target keeps /usr/local/bin and the cgroup
# hierarchy ready; Before=docker.service makes scx's ExecStart fire before dockerd's. This is
# a SOFT order, not a hard gate: Type=simple marks the unit started at fork, so docker is not
# held until flash has attached - in practice flash wins because docker also waits on
# network-online.target (~13-14s in) while flash attaches in ~1s (a few s more if attempt 1
# hits the early-boot ENOMEM and retries), and the try_primary retry below covers any miss.
# (Was After=multi-user.target, which fired the unit ~68s in, mid-storm, behind a
# plymouth-gated target, and left a queued-but-unrun start job.)
After=basic.target
Before=docker.service
ConditionPathIsDirectory=/sys/kernel/sched_ext

[Service]
Type=simple
Restart=on-failure
RestartSec=10
ExecStart=/usr/local/sbin/scx-servermax-start.sh

[Install]
WantedBy=multi-user.target
