# systemd/nuc16pro-servermax-cpupower.service
[Unit]
Description=NUC 16 Pro ServerMax CPU full performance policy (Panther Lake P/E/LP-E)
After=multi-user.target

[Service]
Type=oneshot
ExecStart=/bin/sh -c 'for g in /sys/devices/system/cpu/cpu*/cpufreq/scaling_governor; do [ -w "$g" ] && echo performance > "$g" || true; done'
ExecStart=/bin/sh -c 'for e in /sys/devices/system/cpu/cpu*/cpufreq/energy_performance_preference; do [ -w "$e" ] && echo performance > "$e" || true; done'
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
